# user/gstreamer/gstreamer1.0/etc/profile.d/gst1_0_env.sh
# Setting up environment variable for Gstreamer

GST_PLUGIN_SYSTEM_PATH_1_0=/usr/lib/gstreamer-1.0
GST_PLUGIN_SCANNER_1_0=/usr/lib/gstreamer1.0/gstreamer-1.0/gst-plugin-scanner

export GST_PLUGIN_SYSTEM_PATH_1_0 GST_PLUGIN_SCANNER_1_0
