# vendors/AtmarkTechno/Armadillo-840/etc/udev/rules.d/z70-input-id.rules
SUBSYSTEM!="input", GOTO="input-id_end"

KERNEL=="event*", \
	RUN+="/etc/config/gen-input-run-data.sh %E{ACTION} %E{DEVPATH}"

LABEL="input-id_end"
