# user/udev/udev-105/etc/udev/rules.d/z50-firmware.rules
ACTION=="add", SUBSYSTEM=="firmware", ENV{FIRMWARE}=="*", RUN+="/lib/udev/firmware.sh"
