# vendors/AtmarkTechno/Armadillo-230.Recover/etc/default/root.crontab
*/15 * * * * /etc/snort/send-alert-email
