# vendors/AtmarkTechno/Armadillo-420.WLAN-AWL12/etc/inittab
::sysinit:/etc/init.d/rc

::respawn:/sbin/getty -L 115200 ttymxc1 vt102

::shutdown:/etc/init.d/reboot
::ctrlaltdel:/sbin/reboot
