# user/qt5/etc/profile.d/qt5_env.sh
# Setting up environment variable for Qt5

QT_PLUGIN_PATH=/lib/qt5/plugins
QML2_IMPORT_PATH=/lib/qt5/qml

export QT_PLUGIN_PATH QML2_IMPORT_PATH
