# vendors/AtmarkTechno/Armadillo-240.Recover/etc/udev/rules.d/z10_slideshow.rules
KERNEL=="sd*",DEVPATH=="/block/sd*/sd*[0-9]",BUS=="usb",ACTION=="add",RUN+="/usr/share/slideshow/slideshow_auto_run start"

KERNEL=="sd*",ACTION=="remove",RUN+="/usr/share/slideshow/slideshow_auto_run stop"
