# vendors/AtmarkTechno/Armadillo-810/etc/udev/rules.d/z60-acm-cert.rules
DRIVERS=="acm", ENV{MODULE}=="certification", ACTION=="change", RUN+="/etc/config/acm-cert.sh $ENV{CODEC}"
