# vendors/AtmarkTechno/Armadillo-420/etc/default/rc.local
#!/bin/sh

. /etc/init.d/functions

PATH=/bin:/sbin:/usr/bin:/usr/sbin

#
# Add your commands to run at boot time
#
