# vendors/AtmarkTechno/Armadillo-WLAN.Common/etc/udev/rules.d/z05_aerial.rules
KERNEL=="wlan*", SUBSYSTEM=="net", ACTION=="add", RUN+="/etc/config/aerial-firmware-load.sh %k"
