# vendors/AtmarkTechno/Armadillo-WLAN.Common/etc/udev/rules.d/z05_awl13.rules
KERNEL=="awlan*", SUBSYSTEM=="net", ACTION=="add", RUN+="/etc/config/awl13-firmware-load.sh %k"
